# ROOT_BACKUP/etc/systemd/system/netctl@jio.service.d/profile.conf
[Unit]
Description=Automatically generated profile by wifi-menu
BindsTo=sys-subsystem-net-devices-wlo1.device
After=sys-subsystem-net-devices-wlo1.device
